# systemd/zombie-shot.service
[Unit]
Description=EVE Trader ZombieShot (publish DISPLAY capture to Redis)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=seb
WorkingDirectory=%h/eve-trader
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/%U
Environment=REDIS_URL=redis://127.0.0.1:6379/0
Environment=ZOMBIE_SHOT_REDIS_KEY=system:zombie:screenshot
Environment=ZOMBIE_SHOT_OUT_DIR=/tmp/eve-trader
Environment=ZOMBIE_SHOT_FORMAT=gif
Environment=ZOMBIE_SHOT_GIF_FRAMES=5
Environment=ZOMBIE_SHOT_GIF_FRAME_INTERVAL_MS=800
Environment=ZOMBIE_SHOT_MIN_GIF_BYTES=40960
ExecStart=%h/eve-trader/.venv/bin/python %h/eve-trader/scripts/zombie_shot.py --display :0 --redis-url redis://127.0.0.1:6379/0 --redis-key system:zombie:screenshot --out-dir /tmp/eve-trader --format gif --gif-frames 5 --gif-frame-interval-ms 800 --min-gif-bytes 40960
StandardOutput=journal
StandardError=journal

